# spawnfcgi.ini
[program:spawnfcgi]
command = spawn-fcgi -a 127.0.0.1 -p 9007 -P /var/run/smokeping-fastcgi.pid -n -u nginx -g nginx -f /opt/smokeping/htdocs/smokeping.fcgi
autostart = true
autorestart = true
user = root
